# applications/chef/roles/chef-dk-installer/templates/etc/profile.d/ruby.sh.j2
export PATH=${PATH}:/opt/chefdk/embedded/bin